# openwrt/oonf-dlep-proxy-git/files/dlep.init
#!/bin/sh /etc/rc.common

START=82

. /usr/share/libubox/jshn.sh

start()
{
	service_start /usr/sbin/dlep-proxy --set global.fork=true --load /etc/dlep.conf
}

stop()
{
	service_stop /usr/sbin/dlep-proxy
}
