# openwrt/oonf-dlep-proxy-git/files/dlep.init
#!/bin/sh /etc/rc.common

START=82

. /usr/share/libubox/jshn.sh

start()
{
	service_start /usr/sbin/dlep_proxy --load /etc/dlep.conf --set global.fork=true 
}

stop()
{
	service_stop /usr/sbin/dlep_proxy
}
